# systemd/user/hypridle.service.d/override.conf
[Unit]
PartOf=hyprland.target

[Install]
WantedBy=
WantedBy=hyprland.target
