# systemd/user/hyprland.target
[Unit]
Description=User services specific to Hyprland session
BindsTo=graphical-session.target
After=graphical-session.target

